# renome-cms.service
[Unit]
Description=Renome CMS
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/renome-cms
ExecStart=/root/renome-cms/scripts/run-dev.sh
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=development

[Install]
WantedBy=multi-user.target
